# templates/frigate-archive.timer.j2
[Unit]
Description=Run Frigate Archival Script Daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
